# mkosi.extra/usr/lib/systemd/user/tailscale-systray.service
[Unit]
Description=Tailscale System Tray
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/tailscale systray

[Install]
WantedBy=default.target
