# mkosi.extra/usr/share/factory/etc/profile.d/mise.sh
[[ -n "${BASH_VERSION:-}" ]] || return 0

eval "$(/usr/bin/mise activate bash)"
