# mkosi.extra/usr/share/factory/etc/profile.d/starship.sh
[[ -n "${BASH_VERSION:-}" ]] || return 0

eval "$(/usr/bin/starship init bash)"
