# mkosi.extra/usr/share/factory/etc/profile.d/zoxide.sh
[[ -n "${BASH_VERSION:-}" && $- == *i* ]] || return 0

eval "$(/usr/bin/zoxide init bash)"
