# system_files/etc/profile.d/atuin.sh
[[ -n "${BASH_VERSION:-}" ]] || return 0

eval "$(/usr/bin/atuin init bash)"
