# system_files/etc/profile.d/zoxide.sh
[[ -n "${BASH_VERSION:-}" ]] || return 0

eval "$(/usr/bin/zoxide init bash)"
